# systemd/netviz-exporter.service
[Unit]
Description=netviz: rebuild snapshot.json from Observium DB

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/netviz
EnvironmentFile=%h/workspace/netviz/.env
ExecStart=%h/workspace/netviz/.venv/bin/python -m netviz.exporter
Nice=10
